# supervisord.conf
[supervisord]
nodaemon=true

[program:agent]
autostart=true
autorestart=true
command=%(ENV_WORKDIR)s/control start
startsecs=10
startretries=5
stopwaitsecs=20
